# opengemini_runtime/opengemini-runtime-bot.service
[Unit]
Description=OpenGemini Runtime Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/fallman/.openclaw/workspace/opengemini_runtime
EnvironmentFile=/home/fallman/.openclaw/workspace/opengemini_runtime/.env
ExecStart=/home/fallman/.openclaw/workspace/opengemini_runtime/run_telegram_bot.sh
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
